# Platforms/TNIO AI Bot/Configuration/Systemd/ollama.service
[Unit]
Description=Ollama local model server for lore embeddings
After=network-online.target

[Service]
Type=simple
Environment=OLLAMA_HOST=127.0.0.1:11434
ExecStart=/usr/local/bin/ollama serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
